# roles/openclaw/templates/openclaw-gateway.service.j2
[Unit]
Description=OpenClaw Gateway (v{{ openclaw_version }})
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
ExecStart={{ openclaw_node_bin }} {{ openclaw_gateway_exec_path }} gateway --port {{ openclaw_gateway_port }}
Restart=always
RestartSec=5
RestartPreventExitStatus=78
TimeoutStopSec=30
TimeoutStartSec=30
SuccessExitStatus=0 143
KillMode=control-group
Environment=HOME={{ openclaw_user_home }}
Environment=TMPDIR=/tmp
Environment=NODE_EXTRA_CA_CERTS=/etc/ssl/certs/ca-certificates.crt
Environment=PATH={{ openclaw_node_bin_dir }}:{{ openclaw_user_home }}/.local/bin:{{ openclaw_user_home }}/.npm-global/bin:{{ openclaw_user_home }}/bin:{{ openclaw_user_home }}/.nix-profile/bin:{{ openclaw_user_home }}/.local/share/pnpm:/usr/local/bin:/usr/bin:/bin
Environment=OPENCLAW_GATEWAY_PORT={{ openclaw_gateway_port }}
Environment=OPENCLAW_SYSTEMD_UNIT={{ openclaw_gateway_service_name }}
Environment="OPENCLAW_WINDOWS_TASK_NAME=OpenClaw Gateway"
Environment=OPENCLAW_SERVICE_MARKER=openclaw
Environment=OPENCLAW_SERVICE_KIND=gateway
Environment=OPENCLAW_SERVICE_VERSION={{ openclaw_version }}

[Install]
WantedBy=default.target
